# files/etc/uci-defaults/resize
#!/bin/sh

resize2fs /dev/mmcblk1p2

rm -f /tmp/luci-indexcache
exit 0